# rootfs/php/sv/40-supervisor/studs/laravel-worker.conf
[program:$TMP_LV_NAME-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --queue=$TMP_LV_NAME --sleep=3 --tries=3 --max-time=3600
directory=/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=ubuntu
numprocs=$TMP_LV_WORKERS
redirect_stderr = true
stdout_logfile_maxbytes=0
stopwaitsecs=3600
priority=300
stdout_logfile=${TMP_LV_LOGS}

