# rootfs/runit/supervisor/studs/octane.ini
[program:octane]
process_name=octane
command=php artisan octane:start --workers=${OCTANE_WORKERS} --task-workers=${OCTANE_TASK_WORKERS} --max-requests=${OCTANE_MAX_REQUESTS}
directory=/home/docker/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=docker
numprocs=1
redirect_stderr = true
stdout_logfile_maxbytes=0
stopwaitsecs=3600
priority=900
stdout_logfile=/dev/fd/1
