# rootfs/runit/supervisor/studs/reverb.ini
[program:reverb]
process_name=reverb
command=php artisan reverb:start --host="0.0.0.0" --port=8080
directory=/home/docker/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=docker
numprocs=1
redirect_stderr = true
stdout_logfile_maxbytes=0
stopwaitsecs=3600
priority=900
stdout_logfile=/dev/fd/1
