# rootfs/runit/supervisor/studs/schedule.ini
[program:schedule]
process_name=schedule
command=php artisan schedule:work
directory=/home/docker/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=docker
numprocs=1
redirect_stderr = true
stdout_logfile_maxbytes=0
stopwaitsecs=3600
priority=900
stdout_logfile=/dev/fd/1
